# fixitplease.service
[Unit]
Description=Fix it Please - Static Web App
After=network.target

[Service]
Type=simple
User=sk
WorkingDirectory=/home/sk/dev/tenant-toolbox-mobile
ExecStart=/usr/bin/npx serve dist -l 5173 -s
Restart=always
RestartSec=5
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target
